# service/device-monitor.service
[Unit]
Description=Linux Device Monitor Agent
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
ExecStart=/usr/local/bin/device-monitor -c /etc/device-monitor/config.json
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
